# spring1.service
[Unit]
Description=Spring1 Application
After=syslog.target network.target

[Service]
Type=simple
User=arturs
Group=arturs
WorkingDirectory=/home/arturs/serveris
ExecStart=/usr/bin/java -jar -Dspring.profiles.active=prod -Dserver.port=8080 /home/arturs/serveris/spring1-0.0.1-SNAPSHOT.jar
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spring1

# Security
NoNewPrivileges=true

# Environment variables for production
Environment=DB_URL=jdbc:postgresql://localhost:5432/spelite
Environment=DB_USERNAME=your-db-username
Environment=DB_PASSWORD=your-db-password
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=SERVER_PORT=8080
Environment=LOGGING_LEVEL_ROOT=INFO
Environment=LOGGING_LEVEL_ORG_ARTURS=DEBUG

[Install]
WantedBy=multi-user.target